# crontab
1   */2 *   *   *     git -C ~/racing pull -q --all
*/2 *   *   *   *     ~/genracerot.sh

