# sv/open_new_2/run
#!/bin/sh
SID=open_new_2
cd "$HOME/srv/$SID" || exit 111
exec "./$SID" -glog
