# sv/open_3/run
#!/bin/sh
SID=open3
cd "$HOME/srv/$SID" || exit 111
exec "./$SID" -glog
